# deploy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:api]
command=/usr/bin/dotnet /app/api/Comprendo.Api.dll
; Hereda ConnectionStrings__DefaultConnection, Jwt__*, etc. desde Render (ver start.sh)
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=10

[program:bot]
command=/start-bot.sh
; Hereda TELEGRAM_BOT_TOKEN, GROQ_API_KEY, INTEGRATION_API_KEY desde Render
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=20

[program:frontend]
command=/start-frontend.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=30

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autostart=true
autorestart=true
priority=40
